# ansible/roles/postgresql/templates/backup.crontab
# Periodically refresh backup in S3
3 0 * * * {{postgresql_backup_command}}
