# ansible/roles/postgresql/templates/postgres-supervisor.conf
[program:postgres]
command=/usr/bin/postgres -D {{postgresql_data_directory}} -c config_file={{postgresql_conf_directory}}/postgres.conf
user=postgres
autorestart=true
